# deploy/meeting-copilot-engine.service
[Unit]
Description=Meeting Copilot Engine
After=network.target

[Service]
Type=simple
User=agent4
Group=agency
WorkingDirectory=/opt/agency-workspace/meeting-copilot/engine
ExecStart=/opt/agency-workspace/meeting-copilot/engine/.venv/bin/python main.py
Restart=always
RestartSec=5
EnvironmentFile=/opt/agency-workspace/meeting-copilot/engine/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=copilot-engine

[Install]
WantedBy=multi-user.target
